# dot_config/systemd/user/awww.service.tmpl
[Unit]
Description=AWWW Wallpaper Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=awww-daemon
Restart=on-failure

[Install]
WantedBy=graphical-session.target
