# dot_config/systemd/user/swayidle.service.tmpl
[Unit]
Description=Sway Idle Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=swayidle -w \
          lock 'hyprlock' \
          timeout 300 'loginctl lock-session' \
          timeout 600 'niri msg action power-off-monitors' \
          resume 'niri msg action power-on-monitors' \
          timeout 900 'systemctl suspend' \
          before-sleep 'hyprlock & sleep 1'
Restart=on-failure

[Install]
WantedBy=graphical-session.target
